# scripts/com.thensmt.codex-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thensmt.codex-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/david/Downloads/Claude/NSMT/content-pipeline/.venv/bin/python</string>
        <string>/Users/david/Downloads/Claude/NSMT/content-pipeline/scripts/codex_review.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/david/Downloads/Claude/NSMT/content-pipeline</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/david/Library/Logs/nsmt-codex-review.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/david/Library/Logs/nsmt-codex-review.err</string>
</dict>
</plist>
