# scripts/com.thensmt.codex-rewrite.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thensmt.codex-rewrite</string>

    <!-- Polls #drafts-pre-codex for new Sonnet drafts produced by CI,
         fact-checks + surgically rewrites with Codex GPT-5.5 xhigh, then
         triggers .github/workflows/publish-corrected.yml on CI to do the
         actual admin POST and #recap-pipeline post.

         Schedule: fires at 06:30 local, polls for up to 120 min (covers
         GH Actions cron 10 UTC drafts that land 06:00-08:00 local with
         the usual CI delay; 120 min covers slow runs up to ~08:30).

         Single window only (was 06:30 + 07:30 — could overlap because
         each poll ran 90 min). codex_rewrite.py also takes an exclusive
         fcntl lock on .codex-rewrite.lock so even a manual run during a
         launchd window exits cleanly instead of racing.                 -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/david/Downloads/Claude/NSMT/content-pipeline/.venv/bin/python</string>
        <string>/Users/david/Downloads/Claude/NSMT/content-pipeline/scripts/codex_rewrite.py</string>
        <string>--wait-minutes</string>
        <string>120</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/david/Downloads/Claude/NSMT/content-pipeline</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/david/Library/Logs/nsmt-codex-rewrite.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/david/Library/Logs/nsmt-codex-rewrite.err</string>
</dict>
</plist>
